# provisioning/upstart/docker-foomo-apps.conf
description "docker foomo apps"
author "Me"
start on filesystem and started docker and vagrant-mounted
stop on runlevel [!2345]
respawn
script
  # Wait for docker to finish starting up first.
  FILE=/var/run/docker.sock
  while [ ! -e $FILE ] ; do
    inotifywait -t 2 -e create $(dirname $FILE)
  done
  /usr/bin/docker start -a foomo-apps
end script
post-stop script
    docker kill foomo-apps
end script